# profile.d/homebin.bash
#
# add ~/bin to the path if it exists

_assure_in_path ~/bin
# if this exists, add it too
_assure_in_path "~/.local/bin"
