# profile.d/logging.bash
#
# set up logging for the startup scripts

source $CONFIG_DIR/scripts/logr.bash
logr start profile
logr "Logging configured"
