# profile.d/logging.bash
#
# set up logging for the startup scripts

source $CONFIG_DIR/lib/logr.bash
logr start profile
logr logging configured
